# crontab/production
* * * * * /var/www/insided/app/yiic.php fundafetch >> /var/www/insided/insided/app/fundafetch.log 2>&1
